# playbooks/files/etc/systemd/system/eer.service
[Unit]
Description=Eastern Edge Robotics Control Software

[Service]
TimeoutStartSec=0
Restart=on-failure
RuntimeDirectory=eer
ExecStart=/opt/eer-{{ version }}/bin/{{ entry_point }} --default /opt/eer-{{ version }}/defaultConfig.yml --config /opt/eer-{{ version }}/config.yml

[Install]
WantedBy=multi-user.target
