# playbooks/files/etc/systemd/system/eer.service
[Unit]
Description=Eastern Edge Robotics Control Software

[Service]
TimeoutStartSec=0
Restart=on-failure
RuntimeDirectory=eer
ExecStart=/usr/bin/java -cp "/opt/eer-{{ version }}/libs/*" {{ entry_point }} \
    --default /opt/eer-{{ version }}/defaultConfig.yml \
    --config /home/{{ ansible_user }}/config.yml

[Install]
WantedBy=multi-user.target
