# playbooks/files/etc/systemd/system/eer.service
[Unit]
Description=Eastern Edge Robotics Control Software

[Service]
User={{ ansible_user }}
TimeoutStartSec=0
Restart=on-failure
ExecStart=/opt/eer-{{ version }}/bin/{{ entry_point }} --default /opt/eer-{{ version }}/defaultConfig.yml --config /opt/eer-{{ version }}/config.yml

[Install]
WantedBy=multi-user.target
