# playbooks/files/etc/network/if-up.d/dhcp-server
#!/bin/sh

if [ "$IFACE" = "{{ ssh_interface }}" ]
then
    service isc-dhcp-server restart
fi
